# custom_CustoPiZer/workspace/scripts/files/cron/pioreactor.cron
0 0 */5 * 7 /usr/local/bin/pio run backup_database
0 0 */29 * * pio log -m "Emptying UI's static dir" -l info && rm /var/www/pioreactorui/static/exports/*.zip /var/www/pioreactorui/static/exports/*.csv
